# serv_files/gadget_init.service
[Unit]
Description=Gadget Boot Starter
After=network.target sound.target

[Service]
Type=oneshot
ExecStart=/home/top/AudioBerry/uac2/gadget_init.sh
User=root
Group=root
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
